# searchweb/config.ini
[server-config]
base_dir = /Users/dr/Documents/code/Python/x_search/searchweb
port = 8080
db_path = /Users/dr/Documents/code/Python/x_search/searchweb/db.sql
static_path = /Users/dr/Documents/code/Python/x_search/searchweb/static
keyfile = 
certfile = 
search-root = /tmp/

[loggers]
keys = root

[logger_root]
level = DEBUG
handlers = hand01,hand02

[handlers]
keys = hand01,hand02

[formatters]
keys = simple

[handler_hand01]
class = StreamHandler
level = INFO
formatter = simple
args = (sys.stderr,)

[handler_hand02]
class = handlers.RotatingFileHandler
level = INFO
args = ('/tmp/searchweb.log', 'a', 10*1024*1024, 5)

[program:searchweb]
command = /usr/local/bin/m-server -c /Users/dr/Documents/code/Python/x_search/searchweb/config.ini
stdout_logfile = /tmp/searchweb.log
stderr_logfile = /tmp/searchweb.err.log

[formatter_simple]
format = %(asctime)s | %(levelname)s [%(filename)s:%(lineno)s]: %(message)s  

